# root/etc/services.d/sonarr/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid abc /app/sonarr/bin/Sonarr \
    -nobrowser -data=/config/sonarr
