# root/etc/services.d/prowlarr/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid abc /app/prowlarr/bin/Prowlarr \
    -nobrowser -data=/config/prowlarr
