# root/etc/s6-overlay/s6-rc.d/init-arr-config/up
/etc/s6-overlay/s6-rc.d/init-arr-config/run
